# files/etc/systemd/system/gvisor-auto-update.timer
[Unit]
Description=gvisor-auto-update daily

[Timer]
OnCalendar=daily
AccuracySec=1m
RandomizedDelaySec=2h
Persistent=true

[Install]
WantedBy=timers.target
